# deploy/systemd/omo.service
[Unit]
Description=OMO Boundary Operations Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omo
Group=omo
ExecStart=/usr/local/bin/omo serve --addr 127.0.0.1:8080 --data /var/lib/omo/omo.db --caddy-config /etc/caddy/omo/omo.caddy --sing-box /usr/local/bin/sing-box --sing-box-config /var/lib/omo/sing-box/config.json --backup-dir /var/backups/omo --update-work-dir /var/lib/omo/updates
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/omo /etc/omo /var/log/omo /var/backups/omo /etc/caddy/omo

[Install]
WantedBy=multi-user.target
