# deploy/systemd/omo.service
[Unit]
Description=OMO Boundary Operations Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omo
Group=omo
ExecStart=/usr/local/bin/omo serve --addr 127.0.0.1:8080
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/omo /etc/omo /var/log/omo /var/backups/omo

[Install]
WantedBy=multi-user.target
